# pkg/gentoo/net-analyzer/rackspace-monitoring-agent/files/init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDFILE=/var/run/$RC_SVCNAME.pid

depend() {
    use net
    use logger dns
}

start() {
    ebegin "Starting $RC_SVCNAME"
    start-stop-daemon --background --start --exec /usr/bin/$RC_SVCNAME \
      --make-pidfile --pidfile $PIDFILE
    eend $?
}

stop() {
    ebegin "Stopping $RC_SVCNAME"
    start-stop-daemon --stop --exec /usr/bin/$RC_SVCNAME \
      --pidfile $PIDFILE
    eend $?
}
